# launchd/com.znehraks.cyber-office-v2.ceo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.znehraks.cyber-office-v2.ceo</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>source "__ROOT__/scripts/load-launchd-env.sh"; exec "__ROOT__/scripts/run-discord-ceo.sh"</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__HOME__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CO_LAUNCHD_ENV_FILE</key>
    <string>__ENV_FILE__</string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/ceo.out.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/ceo.err.log</string>
</dict>
</plist>
